# tool/automation/lucide_update.timer
[Unit]
Description=Run lucide update daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

